# crontab
# this is cron job file
0 2 * * * /home/user/project/scripts/backup.sh
