# crontab
0 2 * * * /home/user/project/scripts/backup.sh
